# latest/dockerinit.ini
[program:dockerinit]
command=/etc/supervisor.d/rundocker.sh
