# latest/dockerinit
#!/sbin/openrc-run
 
name="dockerinit"
description="Docker init script"
 
start() {
    ebegin "Starting docker"
    source /docker-lib.sh
    start_docker
    eend $?
}

stop() {
    ebegin "Stopping docker"
    source /docker-lib.sh
    stop_docker
    eend $?
}
