# latest/dockerinit.conf
[program:dockerinit]
command=source /docker-lib.sh && start_docker
